# scripts/crontab
0 * * * * cd /app && /usr/local/bin/python download_submission.py  >> data/download_submission.log 2>&1
* * * * * cd /app && /usr/local/bin/python py_db_test.py  >> data/db_test.log 2>&1

