# scripts/crontab
0 * * * * cd /app && python download_submission.py  >> data/download_submission.log 2>&1

