# scripts/crontab
0 * * * * cd /app && /usr/local/bin/python download_submission.py  >> data/download_submission.log 2>&1
* * * * * cd /app && /usr/local/bin/python scripts/py_db_test.py  >> db_test.log 2>&1

