# scripts/crontab
*/1 * * * * echo test >> echotest.log 2>&1

